# deploy/awp25-collector.service
[Unit]
Description=AWP25 near-real-time reflector telemetry collector
After=local-fs.target

[Service]
Type=simple
User=awp25-monitor
Group=awp25-monitor
WorkingDirectory=/opt/awp25-monitor
EnvironmentFile=/etc/awp25-monitor.env
ExecStart=/usr/bin/python3 -u /opt/awp25-monitor/collector/collector.py
Restart=on-failure
RestartSec=2
StateDirectory=awp25-monitor
StateDirectoryMode=0755
UMask=0022

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
